# supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=npm run serve

[program:nginx]
command=nginx
